# udev/99-keyboard.rules
# 99-keyboard.rules — re-apply keyboard settings whenever a USB keyboard appears
#
# ID_INPUT_KEYBOARD is set by udev's input_id builtin for any device that
# presents keyboard capabilities. Fires on:
#   - USB keyboard physically plugged in
#   - KVM switches that re-enumerate the USB device on port switch
#
# keyboard-udev (root) bridges into the user's systemd session and starts
# keyboard-reset.service, which calls ~/bin/keyboard-settings.

SUBSYSTEM=="input", ACTION=="add", \
    ENV{ID_INPUT_KEYBOARD}=="1", \
    ENV{ID_BUS}=="usb", \
    RUN+="/usr/local/bin/keyboard-udev"
