# home/.config/systemd/user/keyboard-reset.service
[Unit]
Description=Re-apply keyboard and mouse settings
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
Environment=DISPLAY=:0
ExecStart=%h/bin/keyboard-settings
# Exit code 0 or 1 are both success (1 = X not ready, logged and skipped)
SuccessExitStatus=0 1

[Install]
WantedBy=graphical-session.target
